# .launchagents/ai.thejarvis.openclaw.daily-status.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.thejarvis.openclaw.daily-status</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jarvis/Mac-Mini-Obsidian-Vault/1. openclaw/.scripts/daily-status-report.sh</string>
    </array>
    
    <!-- Run daily at 08:00 (8 AM) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/jarvis/Mac-Mini-Obsidian-Vault/1. openclaw/daily-status.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/jarvis/Mac-Mini-Obsidian-Vault/1. openclaw/daily-status.error</string>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
